# luci-cjdns/files/etc/uci-defaults/luci-theme-meshbox
uci set luci.themes.Meshbox=/luci-static/meshbox; uci commit luci
