# cjdns/files/etc/init.d/cjdns
#!/bin/sh /etc/rc.common

START=90
STOP=85

USE_PROCD=1

start_service()
{
	[ -f /etc/uci-defaults/cjdns ] && ( . /etc/uci-defaults/cjdns )

	procd_open_instance
	procd_set_param respawn
	procd_set_param command /bin/ash -c "cjdrouteconf get | cjdroute"
	procd_close_instance

	procd_open_instance
	procd_set_param respawn
	procd_set_param command cjdrouteconf listen
	procd_close_instance
}

stop_service()
{
	procd_kill cjdns
	# FIXME: procd_kill only kills the shell, and cjdroute itself stays open :/
	killall cjdroute || true
}

reload_service()
{
	date >> /tmp/reload_service
	ubus send cjdns.reload '{}'
}

# FIXME: this is broken, i.e. it doesn't call reload_service
service_triggers()
{
	date >> /tmp/service_triggers
	procd_add_reload_trigger cjdns
}
